# scripts/icar_web_gateway.service
[Unit]
Description=iCar Web Control Gateway
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=jetson
ExecStartPre=/bin/bash -c 'cd /home/jetson/icar-ros2-patrol'
ExecStart=/usr/bin/python3 /home/jetson/icar-ros2-patrol/app/web_gateway.py
WorkingDirectory=/home/jetson/icar-ros2-patrol
Environment=ICAR_BRIDGE_HOST=127.0.0.1
Environment=ICAR_BRIDGE_PORT=6501
Environment=ICAR_CAMERA_DEVICE=/dev/camera_depth
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
