# scripts/icar_web_gateway.service
[Unit]
Description=iCar Web Control Gateway
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStartPre=/bin/bash -c 'if [ -d /home/jetson/icar-deploy/current ]; then cd /home/jetson/icar-deploy/current; else cd /home/jetson/icar-ros2-patrol; fi'
ExecStart=/bin/bash -c 'if [ -f /home/jetson/icar-deploy/current/app/web_gateway.py ]; then DIR=/home/jetson/icar-deploy/current; else DIR=/home/jetson/icar-ros2-patrol; fi; exec /usr/bin/python3 "$DIR/app/web_gateway.py"'
WorkingDirectory=/home/jetson/icar-ros2-patrol
Environment=ICAR_BRIDGE_HOST=127.0.0.1
Environment=ICAR_BRIDGE_PORT=6501
Environment=ICAR_CAMERA_DEVICE=/dev/camera_depth
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
